# t8-monitor.service
[Unit]
Description=T8 Delays Monitor - Sydney Airport Line Service Disruption Monitor
After=network.target ollama.service
Wants=ollama.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/t8-monitor
Environment=PATH=/home/pi/t8-monitor/venv/bin
ExecStart=/home/pi/t8-monitor/venv/bin/python monitor_t8_delays_polling.py
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=t8-monitor

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/t8-monitor

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
